# cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Hermes Dashboard
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudflared tunnel run hermes-dash
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
